# Hakurei/hakurei-dind/supervisord.conf
[supervisord]
nodaemon=true

[program:yakumo-instance]
directory=/opt/yakumo-instance/
command=run_server

[program:wrapdocker]
command=/usr/local/bin/wrapdocker

