# Hakurei/hakurei-dind/supervisord.conf
[supervisord]
nodaemon=true

[program:wrapdocker]
command=/usr/local/bin/wrapdocker

